# assets/vpn-enrollment-portal.service
[Unit]
Description=OpenVPN LAN Party enrollment portal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vpnportal
Group=vpnportal
ExecStart=/usr/bin/python3 /usr/local/libexec/vpn-enrollment-portal.py --config /etc/openvpn-lan-party/enrollment-portal.json
Restart=on-failure
UMask=0077
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
ReadWritePaths=/var/lib/openvpn-lan-party/enrollment-portal /run/openvpn-lan-party
RuntimeDirectory=openvpn-lan-party
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
