# assets/vpn-profile-acme-renew.timer
[Unit]
Description=Check the enrollment portal IP certificate twice daily

[Timer]
OnBootSec=15min
OnUnitActiveSec=12h
RandomizedDelaySec=30min
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
